# system/same-eink.path
[Unit]
Description=Watch last_message.json for e-ink updates

[Path]
PathChanged=/home/noaa/sameasy/runtime/last_message.json

[Install]
WantedBy=multi-user.target